# rc.local
  alsactl init
  amixer set Mic 100% cap
